# recipes-core/crash-handler/files/configure-crash-handler.service
[Unit]
Description=Configure Crash Handler

[Service]
Type=oneshot
ExecStart=/usr/sbin/crash_handler --install

[Install]
WantedBy=multi-user.target